# supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=/dind.sh --storage-driver=overlay --group=dockremap

[program:gocd-agent]
command=/gocd-agent.sh